# systemd/synthgen-api.service
[Unit]
Description=SynthGen FastAPI Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/Synteticdatagen
Environment="PATH=/var/www/Synteticdatagen/venv/bin"
ExecStart=/var/www/Synteticdatagen/venv/bin/uvicorn api.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
